# .config/systemd/user/wallpaper.timer
[Unit]
Description=Change wallpapers every 30 minutes

[Timer]
OnActiveSec=30min
OnUnitActiveSec=30min
