# .config/systemd/user/wallpaper.timer
[Unit]
Description=Change wallpapers every 10 minutes

[Timer]
OnActiveSec=10min
OnUnitActiveSec=10min
